# conf/consul.conf
description "Consul server process"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

exec consul agent -server -bootstrap-expect 1 -data-dir /opt/majordomus/majord-data/consul/data -config-dir /opt/majordomus/majord-data/consul/conf